# build/playbooks/files/consul_template.conf
[program:consul-template]
command=/usr/local/bin/consul-template -config /etc/consul-template/conf.d/
autostart=true
autorestart=true
stderr_logfile=/var/log/consul-template.err.log
stdout_logfile=/var/log/consul-template.out.log